# systemd/aurelia-mayor.service
[Unit]
Description=Aurelia Mayor Agent Daemon
After=network.target aurelia-server.service
Requires=aurelia-server.service

[Service]
User=mayor
Group=mayor
ExecStart=/bin/bash /home/zuzu/Code/aurelia/scripts/agent_daemon.sh mayor http://localhost:8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aurelia-mayor

[Install]
WantedBy=multi-user.target
