# systemd/aurelia-janitor.service
[Unit]
Description=Aurelia Janitor Agent Daemon
After=network.target aurelia-server.service
Requires=aurelia-server.service

[Service]
User=janitor
Group=janitor
ExecStart=/bin/bash /home/zuzu/Code/aurelia/scripts/agent_daemon.sh janitor http://localhost:8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aurelia-janitor

[Install]
WantedBy=multi-user.target
