# systemd/aurelia-personal.service
[Unit]
Description=Aurelia Personal Agent Daemon
After=network.target aurelia-server.service
Requires=aurelia-server.service

[Service]
User=personal
Group=personal
ExecStart=/bin/bash /home/zuzu/Code/aurelia/scripts/agent_daemon.sh personal http://localhost:8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aurelia-personal

[Install]
WantedBy=multi-user.target
